# systemd/jupyter.service
[Unit]
Description=JupyterLab (localhost only)
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
Restart=always
RestartSec=5
User=%i
WorkingDirectory=/home/%i/work
ExecStart=/opt/jupyter-venv/bin/jupyter lab --ip=127.0.0.1 --port=8888 --no-browser --ServerApp.token=${JUPYTER_TOKEN}

[Install]
WantedBy=multi-user.target
